# provisioning/systemd/rig-manager.service
[Unit]
Description=Rig Control Plane
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/rig-manager --config /etc/rig-manager/production.toml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
